# bub_events/src/skills/bub-events/templates/cron.timer
# Bub Cron Timer Template
# 
# Install:
#   sudo cp bub-cron-@NAME@.timer /etc/systemd/system/
#   sudo cp bub-cron-@NAME@.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable bub-cron-@NAME@.timer
#   sudo systemctl start bub-cron-@NAME@.timer
#
# Variables:
#   @NAME@      - Job name (e.g., backup-check)
#   @SCHEDULE@  - systemd calendar spec (e.g., daily, *-*-* 03:00:00)
#   @CONTENT@   - Event content message
#   @TOPIC@     - Topic identifier for prompt loading

[Unit]
Description=Bub cron timer for @NAME@

[Timer]
OnCalendar=@SCHEDULE@
Persistent=true

[Install]
WantedBy=timers.target
